# docs/ansible/roles/frr/templates/etc/frr-routing.service.j2
[Unit]
Description=push configuration data stored in sqlite db into frr conf file job

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/opt/netfoundry/nfroute/bin/frr_routing -t 60"
