# quickstarts/docs/ansible/roles/frr/templates/etc/frr-routing.timer.j2
[Unit]
Description=frr-routing

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s
AccuracySec=1s

[Install]
WantedBy=timers.target
