# quickstarts/docs/ansible/roles/frr/templates/etc/mop-routing.service.j2
[Unit]
Description=mop declared configuration pushed into sqlite db job

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/opt/netfoundry/nfroute/bin/mop_routing"
